# crontab
* * * * * /home/pi/develop/bme280/entry-point.sh
