# crontab
* * * * * /bin/bash /home/pi/develop/indoor_air_condition_measure/entry-point.sh
